# projects-app/core/src/main/dist/systemctlctl/java-gradle-app-skeleton.service
[Unit]
Description=My Great Application
Documentation=https://app.java-gradle-app-skeleton.org
After=syslog.target network.target

[Service]
Type=simple
User=root
Group=root
Environment="JAVA_HOME=/opt/java-home-17"
ExecStart=/path/to/java-gradle-app-skeleton/startup.sh
KillSignal=15
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
